# spacer/supervisor/autolapse.conf
[program:autolapse]
command=/home/micah/src/broadcasting/spacer/lapser/autolapse.py
autostart=true
autorestart=true
