# spacer/supervisor/vfd.conf
[program:vfd]
command=/home/micah/vfd/dashboard.py
autostart=true
autorestart=true
