# spacer/supervisor/vfd.conf
[program:vfd]
command=/home/micah/src/broadcasting/spacer/vfd/dashboard.py
autostart=true
autorestart=true
